# profile.d/general.sh
export VISUAL=vim
export EDITOR=vim

# Mainly for i3 (check /usr/bin/i3-sensible-terminal).
export TERMINAL=urxvt

export PAGER=less
